# airootfs/etc/skel/.local/share/systemd/user/alarm.service
[Unit]
Description=Check appointments
Wants=alarm.timer
[Service]
ExecStart=/usr/bin/python alarm.py
[Install]
WantedBy=multi-user.target
