# airootfs/etc/skel/.local/share/systemd/user/cal.service
[Unit]
Description=Sync calendar
Wants=cal.timer
[Service]
ExecStart=/bin/bash calsync.sh
[Install]
WantedBy=multi-user.target
