# airootfs/etc/skel/.local/share/systemd/user/mail.service
[Unit]
Description=Sync mail
Wants=mail.timer
[Service]
ExecStart=/bin/bash sync-mail.sh
[Install]
WantedBy=multi-user.target
